# supervisor/conf.d/ai_parse.conf
[program:ai_parse]
directory=/var/www/html/flask_servers/ai_parse/src
command=gunicorn -w 3 -t 7200 -k gevent -b 0.0.0.0:5001 main:app
