# supervisor/conf.d/local_news_parse.conf
[program:local_news_parse]
directory=/var/www/html/flask_servers/local_news_parse/src
command=gunicorn -w 2 -t 3600 -k gevent -b 0.0.0.0:5000 main:app
